# cron.daily/apf-firewall
#! /bin/sh

. /etc/default/apf-firewall
if [ "x$RUN" != "xyes" ] ; then
	exit 0
fi

/usr/sbin/apf -f >> /dev/null 2>&1
/usr/sbin/apf -s >> /dev/null 2>&1

